# debian/obs-worker.prerm
#!/bin/sh -e

if [ "$1" = "purge" ]; then
    invoke-rc.d obsworker stop || true
    deluser --system --quiet obsworker
    delgroup --system --quiet obsworker
fi

#DEBHELPER#
