# debian/obs-worker.prerm
#!/bin/sh -e

if [ "$1" = "remove" ]; then
    invoke-rc.d obsworker stop || true
fi

#DEBHELPER#
