# dist/goaccess/obs-goaccess-last-month.crontab
0 0 5 * * wwwrun /usr/bin/obs-goaccess-last-month.sh
