# dist/obsapisetup
#! /bin/sh
# Copyright (c) 2010, Novell Inc.
#
# Author: adrian@suse.de
#
# /etc/init.d/obsapisetup
#   and its symbolic  link
# /usr/sbin/rcobsapisetup
#
### BEGIN INIT INFO
# Provides:          obsapisetup
# Start-Before:      lighttpd
# Should-Start:      obsstoragesetup
# Should-Stop:       $none
# Required-Start:    mysql
# Required-Stop:     $null
# Default-Start:     3 5
# Default-Stop:      0 1 2 4 6
# Description:       Initialize and update api database, only used in OBS Appliance
### END INIT INFO

. /etc/rc.status

# Determine the base and follow a runlevel link name.
base=${0##*/}
link=${base#*[SK][0-9][0-9]}

apidir=/srv/www/obs/api
webuidir=/srv/www/obs/webui

# make parsed output predictable 
export LC_ALL=C

FQHOSTNAME=`hostname -f`


if [ "$?" != "0" ]; then
  # Fallback to IP of the VM/host
  FQHOSTNAME=`ip addr | sed -n 's,.*inet \(.*\)/.* brd.*,\1,p' | grep -v ^127. | head -n 1`
  if [ "$?" != "0" -o "$FQHOSTNAME" = "" ]; then
    echo "    Can't determine hostname or IP - Network setup failed!"
    echo "    Check if networking is up and dhcp is working!"
    echo "    Using 'localhost' as FQHOSTNAME."
    FQHOSTNAME="localhost"
  fi
fi
. /etc/sysconfig/obs-server

if [ "$OBS_API_AUTOSETUP" != "yes" ]; then
   echo "OBS API Autosetup is not enabled in sysconfig, skipping!"
   exit 0
fi

rc_reset
case "$1" in
	start)
                if [ ! -d /obs/MySQL ]; then
			mkdir -p /obs/MySQL
		fi
		echo -n "Adjust webui configuration for api"
                # it would be better to solve this in lighttpd.conf to accept also localhost:81,
                # but no idea how to solve this
                sed -i 's,^FRONTEND_HOST = .*,FRONTEND_HOST = "'"$FQHOSTNAME"'",' \
		    /srv/www/obs/webui/config/environments/production.rb
                sed -i 's,^DOWNLOAD_URL = .*,DOWNLOAD_URL = \"http://'"$FQHOSTNAME:82"'\",' \
		    /srv/www/obs/webui/config/environments/production.rb
                sed -i -e /webui_url:/d -e /webui_host:/d -e /allow_anonymous:/d \
		    /srv/www/obs/api/config/options.yml
                echo "# the following lines are always written by obsapisetup!" >> /srv/www/obs/api/config/options.yml
                echo "allow_anonymous: true" >> /srv/www/obs/api/config/options.yml
		if [ "$FQHOSTNAME" == "localhost" ] ; then
			echo "webui_host: localhost" >> /srv/www/obs/api/config/options.yml
		else
			echo "webui_host: `ip addr | sed -n 's,.*inet \(.*\)/.* brd.*,\1,p' | grep -v ^127. | head -n 1`" >> /srv/www/obs/api/config/options.yml
		fi
                echo "webui_url: http://$FQHOSTNAME:80" >> /srv/www/obs/api/config/options.yml
		echo
		OBSVERSION=`rpm -q --qf '%{VERSION}' obs-server`
		OS=`head -n 1 /etc/SuSE-release`
                if [ ! -e /obs/MySQL/api_production ]; then
			echo -n "Initialize OBS api database (first time only)"
			mysqladmin -u root create api_production
			mysqladmin -u root password "opensuse"
        		cd $apidir
			RAILS_ENV=production rake db:setup >> $apidir/log/db_migrate.log
                        ( umask 0077; dd if=/dev/urandom bs=256 count=1 2>/dev/null |sha256sum| cut -d\  -f 1 >config/secret.key )
			chown lighttpd.lighttpd config/secret.key
		else
			echo
			echo -n "Migrate OBS api database"
			cd $apidir
			RAILS_ENV=production rake db:migrate >> $apidir/log/db_migrate.log
			echo
		fi
                if [ ! -e /obs/MySQL/webui_production ]; then
			echo -n "Initialize OBS webui database (first time only)"
			mysqladmin -u root create webui_production
        		cd $webuidir
			RAILS_ENV=production rake db:setup >> $webuidir/log/db_migrate.log
                        ( umask 0077; dd if=/dev/urandom bs=256 count=1 2>/dev/null |sha256sum| cut -d\  -f 1 >config/secret.key)
			chown lighttpd.lighttpd config/secret.key
		else
			echo
			echo -n "Migrate OBS webui database"
			cd $webuidir
			RAILS_ENV=production rake db:migrate >> $webuidir/log/db_migrate.log
			echo
		fi

                if [ ! -e /etc/ssl/private/lighttpd.key -a ! -e /etc/ssl/certs/lighttpd.crt -a ! -e /etc/ssl/private/lighttpd.pem ]; then
                        if [ ! -e /OBS/certs/lighttpd.key ]; then
                               # setup ssl certificates (NOT protected with a passphrase)
                               echo "Creating SSL certificate for the server..."
                               install -d -m 0700 /OBS/certs
                               openssl genrsa -out /OBS/certs/lighttpd.key 1024
                               openssl req -new -key /OBS/certs/lighttpd.key -out /OBS/certs/lighttpd.csr
                               openssl x509 -req -days 365 -in /OBS/certs/lighttpd.csr -signkey /OBS/certs/lighttpd.key -out /OBS/certs/lighttpd.crt
                               cat /OBS/certs/lighttpd.key /OBS/certs/lighttpd.crt > /OBS/certs/lighttpd.pem
                        fi
                        # use certificates from /OBS (LVM ?) partition
                        ln -s /OBS/certs/lighttpd.key /etc/ssl/private/lighttpd.key
                        ln -s /OBS/certs/lighttpd.crt /etc/ssl/certs/lighttpd.crt
                        ln -s /OBS/certs/lighttpd.pem /etc/ssl/private/lighttpd.pem
                fi

                cd $webuidir
		chown -R lighttpd.lighttpd $webuidir/log
                cd $apidir
		chown -R lighttpd.lighttpd $apidir/log

		cat > /etc/issue <<EOF
Welcome to openSUSE Build Service(OBS) Appliance $OBSVERSION
based on $OS

EOF
		if ! grep -q "^our \$sign =" /usr/lib/obs/server/BSConfig.pm ; then
			cat >> /etc/issue <<EOF

  WARNING: **** Package signing is disabled, maybe due to lack of hardware number generator ****

EOF
		fi

		if [ -n "$FQHOSTNAME" ]; then
			cat >> /etc/issue <<EOF

  Connect to the web interface via:     http://$FQHOSTNAME
  Connect to the api interface via:     http://$FQHOSTNAME:81
  Browse the build packages via:        http://$FQHOSTNAME:82

 * "Admin" user password is "opensuse" by default.
 * Connect to the web interface now to finish the OBS setup.

More informations about this appliance are available here:

 http://en.opensuse.org/Build_Service/OBS-Appliance

                                  Greetings from the openSUSE Build Service Team

EOF
		else
			echo '**********************************************' >> /etc/issue
			echo '**           NETWORK SETUP FAILED           **' >> /etc/issue
			echo '**                                          **' >> /etc/issue
			echo '** OBS is not usable                        **' >> /etc/issue
			echo '** A working DHCP and DNS server in network **' >> /etc/issue
			echo '** is required!                             **' >> /etc/issue
			echo '**********************************************' >> /etc/issue
		fi
		rc_status -v
	;;
	stop)
                # nothing to do
		rc_status -v
	;;
	restart)
                # nothing to do
		rc_status
	;;
	try-restart)
                # nothing to do
		rc_status
	;;
	reload)
                # nothing to do
		rc_status
	;;
	status)
		# nothing to do
		rc_status -v
	;;
	*)
		echo "Usage: $0 {start|stop|status|try-restart|restart|reload}"
		exit 1
	;;
esac
rc_exit
