# dist/obswarden
#! /bin/sh
# Copyright (c) 2007, Novell Inc.
#
# Author: adrian@suse.de
#         Martin Mohring <martin.mohring@5etech.eu>
#         Susanne Oberhauser  <froh@novell.com>
#
# /etc/init.d/obswarden
#   and its symbolic  link
# /usr/sbin/rcobswarden
#
### BEGIN INIT INFO
# Provides:          obswarden
# Required-Start:    $time $syslog obsrepserver
# Required-Stop:     obssrcserver $null
# Should-Start:      obsstoragesetup
# Should-Stop:       $none
# Default-Start:     3 5
# Default-Stop:      0 1 2 4 6
# Short-Description: OBS service warden
# Description:       open build service warden, monitors the workers
### END INIT INFO

. /etc/rc.status

. /etc/sysconfig/obs-server

# Determine the base and follow a runlevel link name.
base=${0##*/}
link=${base#*[SK][0-9][0-9]}

if [ -z "$OBS_RUN_DIR" ]; then
        OBS_RUN_DIR="/srv/obs/run"
fi

if [ -z "$OBS_LOG_DIR" ]; then
        OBS_LOG_DIR="/srv/obs/log"
fi

if [ -n "$OBS_BACKENDCODE_DIR" ]; then
	obsdir="$OBS_BACKENDCODE_DIR"
else
	obsdir=/usr/lib/obs/server
fi
rundir="$OBS_RUN_DIR"
logdir="$OBS_LOG_DIR"

rc_reset
case "$1" in
	start)
		echo -n "Initializing obswarden"
                # FIXME: not nice, this should receive a proper daemon handling, 
                #        including real logging, pid file and startproc
		startproc -f -l "$logdir"/warden.log \
		    "$obsdir"/bs_warden
		rc_status -v
	;;
	stop)
		echo -n "Shutting down obswarden"
		"$obsdir"/bs_warden --stop
		rc_status -v
	;;
	restart)
		"$obsdir"/bs_warden --restart
		rc_status
	;;
	try-restart|reload)
		$0 status
		if test $? = 0; then
			"$obsdir"/bs_warden --restart
		else
			rc_reset        # Not running is not a failure.
		fi
		# Remember status and be quiet
		rc_status
	;;
	status)
		echo -n "Checking for obswarden: "
		checkproc "$obsdir"/bs_warden
		rc_status -v
	;;
	*)
		echo "Usage: $0 {start|stop|status|try-restart|restart|reload}"
		exit 1
	;;
esac
rc_exit
