# dist/systemd/obs-sphinx.service
[Unit]
Description = Open Build Service Sphinx Search Daemon
BindsTo = obs-api-support.target
Conflicts = searchd.service

[Service]
Environment = "RAILS_ENV=production"
User = wwwrun
Group = www
WorkingDirectory = /srv/www/obs/api
ExecStart = /usr/bin/bundle.ruby2.5 exec rails sphinx:start
ExecStop = /usr/bin/bundle.ruby2.5 exec rails ts:stop
Type = forking
PIDFile = /srv/www/obs/api/log/production.sphinx.pid

[Install]
WantedBy = obs-api-support.target
