# dist/systemd/obsdodup.service
[Unit]
Description=OBS dodup, updates download on demand metadata
After=network-online.target
After=obsapisetup.service

[Service]
OOMPolicy=continue
EnvironmentFile=/etc/sysconfig/obs-server
ExecStart=/usr/lib/obs/server/bs_dodup --logfile dodup.log
ExecStop=/usr/lib/obs/server/bs_dodup --stop
ExecReload=/usr/lib/obs/server/bs_dodup --restart

[Install]
WantedBy=multi-user.target
