# dist/systemd/obsgetbinariesproxy.service
[Unit]
Description=OBS binaries proxy cache
After=network.target obsstoragesetup.service obsstoragesetup.service

[Service]
OOMPolicy=continue
LimitNOFILE=16384
EnvironmentFile=/etc/sysconfig/obs-server
ExecStart=/usr/lib/obs/server/bs_getbinariesproxy --logfile getbinariesproxy.log
ExecStop=/usr/lib/obs/server/bs_getbinariesproxy --stop
ExecReload=/usr/lib/obs/server/bs_getbinariesproxy --restart

[Install]
WantedBy=multi-user.target
