# dist/systemd/obspublisher.service
[Unit]
Description=OBS repository publisher
After=network-online.target

[Service]
OOMPolicy=continue
EnvironmentFile=/etc/sysconfig/obs-server
ExecStart=/usr/lib/obs/server/bs_publish --logfile publisher.log
ExecStop=/usr/lib/obs/server/bs_publish --stop
ExecReload=/usr/lib/obs/server/bs_publish --restart

[Install]
WantedBy=multi-user.target
