# dist/systemd/obspublisher.service
[Unit]
Description=OBS repository publisher
After=network.target

[Service]
EnvironmentFile=/etc/sysconfig/obs-server
ExecStart=/usr/lib/obs/server/bs_publish
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
