# dist/systemd/obsrepserver.service
[Unit]
Description=OBS repository server
After=network-online.target obsstoragesetup.service

[Service]
OOMPolicy=continue
EnvironmentFile=/etc/sysconfig/obs-server
ExecStart=/usr/lib/obs/server/bs_repserver --logfile rep_server.log
ExecStop=/usr/lib/obs/server/bs_repserver --stop
ExecReload=/usr/lib/obs/server/bs_repserver --restart

[Install]
WantedBy=multi-user.target
