# dist/systemd/obsservicedispatch.service
[Unit]
Description=OBS source service dispatcher
After=network-online.target obsstoragesetup.service obsservice.service

[Service]
OOMPolicy=continue
EnvironmentFile=/etc/sysconfig/obs-server
ExecStart=/usr/lib/obs/server/bs_servicedispatch --logfile servicedispatch.log
ExecStop=/usr/lib/obs/server/bs_servicedispatch --stop
ExecReload=/usr/lib/obs/server/bs_servicedispatch --restart

[Install]
WantedBy=multi-user.target
