# dist/systemd/obsworker@.service
[Unit]   
Description=OBS worker, instance %i
After=network-online.target obsstoragesetup.service obssrcserver.service obsrepserver.service xendomains.service
         
[Service]
OOMPolicy=continue
EnvironmentFile=/srv/obs/run/worker/%i/parameters.txt
ExecStart=nice -n $OBS_NICE /srv/obs/run/worker/boot/bs_worker $OBS_WORKER_PARAMETERS
ExecStop=nice -n $OBS_NICE /srv/obs/run/worker/boot/bs_worker --statedir /srv/obs/run/worker/%i --exit
         
[Install]
WantedBy=multi-user.target
